# greetd/etc/systemd/system/greeter-status.service
# ============================================================================
#  Coletor de status do greeter (roda como ROOT)
# ----------------------------------------------------------------------------
#  Gera /run/greeter-status/status.json + assets (wallpaper borrado, GIFs, quote)
#  num loop. O greeter quickshell (usuário `greeter`, sem acesso a docker) só LÊ
#  esse JSON — toda a coleta privilegiada acontece aqui.
# ============================================================================
[Unit]
Description=Greeter live status collector (docker/temps/quote/wallpaper)
After=docker.service
Wants=docker.service

[Service]
Type=simple
ExecStart=/usr/local/lib/greetd/collect-status.sh
RuntimeDirectory=greeter-status
RuntimeDirectoryMode=0755
Restart=always
RestartSec=2
Nice=10

[Install]
WantedBy=multi-user.target
