# caddy/etc/systemd/system/caddy.service.d/env.conf
[Service]
# Carrega as variáveis do .env dos dotfiles no ambiente do processo do Caddy.
# O systemd (PID 1, como root) lê este arquivo antes de aplicar o sandbox,
# então o ProtectHome=true da unit não impede a leitura em /home.
EnvironmentFile=-/home/v1cferr/dotfiles/.env
