# caddy/etc/systemd/system/caddy.service.d/exec.conf
[Service]
# Aponta a unit para o binário CUSTOM em /usr/local/bin/caddy, buildado com
# xcaddy + o plugin dns.providers.cloudflare (necessário para o DNS-01 challenge
# do certificado curinga *.v1cferr.dev). O pacote `caddy` do pacman fica como
# está em /usr/bin/caddy — mas NÃO é usado: um `pacman -Syu` pode atualizá-lo à
# vontade que não derruba o proxy. Para subir a versão do Caddy custom, rode
# scripts/caddy/build.sh (não o pacman).
#
# A sintaxe "Exec*=" vazio antes de reatribuir é exigida pelo systemd para
# substituir o valor herdado da unit original (senão ele ACUMULA os comandos).
ExecStartPre=
ExecStartPre=/usr/local/bin/caddy validate --config /etc/caddy/Caddyfile
ExecStart=
ExecStart=/usr/local/bin/caddy run --config /etc/caddy/Caddyfile
ExecReload=
ExecReload=/usr/local/bin/caddy reload --config /etc/caddy/Caddyfile --force
