# scripts/secrets/dotfiles-secrets-backup.timer
[Unit]
Description=Atualiza o backup criptografado dos segredos (diário)

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
