# scripts/secrets/dotfiles-secrets-backup.service
[Unit]
Description=Backup criptografado dos segredos (SSH/GPG/.env/tokens)

[Service]
Type=oneshot
ExecStart=%h/dotfiles/scripts/secrets/backup.sh
